# data-pipeline/com.async.finance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.async.finance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zhangyang/Desktop/NYU/vibe-coding/a-share-finance/data-pipeline/run_sync.sh</string>
    </array>

    <!-- 18:00 Beijing time = 10:00 UTC, Mon-Fri -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>1</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zhangyang/.a-share-sync-logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zhangyang/.a-share-sync-logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/zhangyang</string>
    </dict>
</dict>
</plist>
